# install/supervisor/judge-queue.conf
[program:judge-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan queue:work --queue=solutions
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/app/storage/logs/queue/judge-queue.log
stopwaitsecs=3600
